# home/users/beepbopbop/files/autostart/discord.desktop
[Desktop Entry]
Type=Application
Name=Discord
Icon=discord
Exec=discord --start-minimized
X-GNOME-Autostart-enabled=true
