# home/users/beepbopbop/files/autostart/solaar.desktop
[Desktop Entry]
Type=Application
Name=Solaar
Icon=solaar
Exec=solaar -w hide
X-GNOME-Autostart-enabled=true
