# home/users/beepbopbop/files/autostart/steam-silent.desktop
[Desktop Entry]
Type=Application
Name=Steam
Icon=steam
Exec=steam -silent
X-GNOME-Autostart-enabled=true
