# home/users/beepbopbop/files/autostart/vesktop.desktop
[Desktop Entry]
Type=Application
Name=Vesktop
Icon=vesktop
Exec=vesktop --start-minimized
X-GNOME-Autostart-enabled=true
